# usr/src/AIC8800/aic.rules
KERNEL=="sd*", ATTRS{idVendor}=="a69c", ATTRS{idProduct}=="5721",  SYMLINK+="aicudisk", RUN+="/usr/local/bin/usb-eject-delay.sh"
